# lib/template.upstart
# Ubuntu upstart file at /etc/init/{user}.conf
# Cf. https://github.com/mongodb/mongo/blob/master/debian/mongod.upstart

start on runlevel [2345]
stop on runlevel [06]

script
  CMD="{command}"
  DAEMON_USER="{user}"
  SERVICE="/usr/local/lib/quichean/$DAEMON_USER/{service}"
  DAEMON_MGR="start-stop-daemon"
  PIDFILE="/var/run/$DAEMON_USER.pid"

  exec $DAEMON_MGR --start \
    --chuid $DAEMON_USER:$DAEMON_USER \
    --pidfile $PIDFILE \
    --make-pidfile \
    --exec $CMD $SERVICE
end script
